# config/files/profile.d/aliases.sh
# WavyOS System Aliases
# Ensures 'just' always finds the system config, even from random directories
alias just='/usr/bin/just --justfile /etc/justfile --working-directory $HOME'
